# deploy/systemd/telegram-relay.service
[Unit]
Description=Telegram to Discord AI Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/friendly-broccoli-discord
EnvironmentFile=/opt/friendly-broccoli-discord/.env
ExecStart=/opt/friendly-broccoli-discord/.venv/bin/python -m goob_ai.main
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
